# contrib/graph/crontab
# Start the script every 5 minutes
*/5 * * * *		root		/usr/local/sbin/chaosvpn.graph
